# zookeeper.conf

[program:zookeeper]
command=bin/zookeeper-server-start.sh config/zookeeper.properties
;user=nobody
priority=10
startretries=5
redirect_stderr=true
